# systemd/nemea-supervisor.service.in
# It's not recommended to modify this file in-place, because it will be
# overwritten during package upgrades.  If you want to customize, the best
# way is to create a file "/etc/systemd/system/nemea-supervisor.service",
# containing
#   .include /lib/systemd/system/nemea-supervisor.service
#   ...make your changes here...
# For more info about custom unit files, see
# http://fedoraproject.org/wiki/Systemd#How_do_I_customize_a_unit_file.2F_add_a_custom_unit_file.3F
#
# For example, to pass additional options (for instance, a different socket path) to the
# supervisor binary at startup, you need to create a file named
# "/etc/systemd/system/nemea-supervisor.service" containing:
#       .include /lib/systemd/system/nemea-supervisor.service
#       [Service]
#       Environment="SOCKET=/my/absolute/path"

[Unit]
Description=Nemea Supervisor
After=network.target local-fs.target remote-fs.target

[Service]
Type=forking
PIDFile=@defaultsocketdir@/nemea-supervisor.pid
Environment="SOCKET=@defaultsocketdir@/nemea-supervisor.sock"
Environment="CONFIG_TEMPLATE=@defaultconfdir@/supervisor_config_template.xml"
Environment="LOG_PATH=/var/tmp/nemea"
Environment=LANG=C
User=nemead
Group=nemead

ExecStart=@bindir@/supervisor -d -s $SOCKET -T $CONFIG_TEMPLATE -L $LOG_PATH
ExecReload=@bindir@/supervisor_cli -r -s $SOCKET
KillSignal=SIGTERM
PrivateTmp=false
PermissionsStartOnly=true

[Install]
WantedBy=multi-user.target

